# recipes-essentials/ecubee-configure/ecubee-configure_1.0/src/ecubee.service
[Unit]
Description=Ecubee Demo App
Requires=pvr-init.service x.service
After=syslog.target pvr-init.service x.service

[Service]
ExecStart=-/usr/bin/ecubee
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=graphical.target
