# recipes-essentials/ecubee-configure/ecubee-configure_1.0/src/network@.service
[Service]
ExecStart=/sbin/ifup %i
Type=oneshot
RemainAfterExit=yes
#ExecStop=/sbin/ifdown %i
StandardOutput=syslog+console

